# lepidopter-fh/etc/cron.daily/remove_upl_reports
#!/usr/bin/env bash
/opt/ooni/remove-upl-reports.sh
