# lepidopter-fh/etc/cron.daily/update_deck
#!/bin/bash
# Update deck
/opt/ooni/update-deck.sh
