# lepidopter-fh/etc/cron.daily/upload_reports
#!/bin/bash
/opt/ooni/upload-reports.sh
