# lepidopter-fh/etc/cron.weekly/remove_inc_reports
#!/bin/bash
/opt/ooni/remove-inc-reports.sh
