# lepidopter-fh/etc/cron.weekly/update_ooniprobe
#!/bin/bash
# Update ooniprobe
/opt/ooni/update-ooniprobe.sh
