# lepidopter-fh/etc/cron.weekly/update_ooniprobe_deck
#!/bin/bash
# Update ooniprobe and deck
/opt/ooni/update-ooniprobe.sh
/opt/ooni/update-deck.sh
