# lepidopter-fh/etc/network/if-up.d/run_oonideckgen
#!/bin/bash
# Update ooni probe on network change

set -e

# Exit for loopback change
if [ "$IFACE" = lo ] ; then
	  exit 0
  fi

/opt/ooni/update-deck.sh

exit 0
